# service/booking.socket
[Unit]
Description=booking gunicorn socket

[Socket]
ListenStream=/run/booking.sock

[Install]
WantedBy=sockets.target